# dotfiles/hermes-gateway.service
[Unit]
Description=Hermes Agent Gateway - Messaging Platform Integration
After=network-online.target ollama.service
Wants=network-online.target ollama.service
StartLimitIntervalSec=0

[Service]
Type=simple
ExecStart=/home/hans/.hermes/hermes-agent/venv/bin/python -m hermes_cli.main gateway run --replace
WorkingDirectory=/home/hans/.hermes/hermes-agent
Environment="PATH=/home/hans/.hermes/hermes-agent/venv/bin:/home/hans/.hermes/hermes-agent/node_modules/.bin:/home/hans/.hermes/node/bin:/home/hans/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="VIRTUAL_ENV=/home/hans/.hermes/hermes-agent/venv"
Environment="HERMES_HOME=/home/hans/.hermes"
Environment="HERMES_ACCEPT_HOOKS=1"
Restart=always
RestartSec=60
RestartMaxDelaySec=300
RestartSteps=5
RestartForceExitStatus=75
KillMode=mixed
KillSignal=SIGTERM
ExecReload=/bin/kill -USR1 $MAINPID
TimeoutStopSec=90
MemoryHigh=1G
MemoryMax=2G
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
