# dotfiles/hermes-mcp.service
[Unit]
Description=HVE Hermes MCP Server
Documentation=https://github.com/humanvalueexchange/hermes-cfo
After=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
WorkingDirectory=%h/hermes-cfo
# Kill any orphaned process holding port 8765 before starting (prevents crash loop)
ExecStartPre=-/usr/bin/fuser -k 8765/tcp
ExecStart=%h/.hermes-mcp-venv/bin/python %h/hermes-cfo/mcp/server.py
EnvironmentFile=-%h/.hermes-mcp.env
Environment="HVE_MCP_HOST=0.0.0.0"
Environment="HVE_MCP_PORT=8765"
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hermes-mcp

[Install]
WantedBy=default.target
