# dotfiles/hermes-data-refresh.service
[Unit]
Description=Hermes — Daily BTC Market Data Refresh
After=network-online.target ollama.service
Wants=network-online.target

[Service]
Type=oneshot
User=hans
WorkingDirectory=/home/hans/freqtrade
# Wait up to 2 minutes for Binance to be reachable (handles post-boot network delays)
ExecStartPre=/bin/sh -c 'for i in $(seq 1 12); do curl -sf --max-time 5 https://api.binance.com/api/v3/ping && exit 0; echo "Waiting for Binance ($i/12)..."; sleep 10; done; echo "Binance unreachable after 2min — aborting"; exit 1'
ExecStart=/home/hans/freqtrade/.venv/bin/freqtrade download-data \
  --exchange binance \
  --pairs BTC/USDT \
  --timeframes 1m 5m 15m 1h 4h 1d \
  --days 7 \
  --datadir /home/hans/freqtrade/user_data/data
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
