# dotfiles/hermes-telegram-log.service
[Unit]
Description=Hermes Telegram Live Log Watcher
Documentation=https://github.com/humanvalueexchange/hermes-cfo

Wants=hermes-gateway.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 /home/hans/hermes-cfo/scripts/telegram-live-log.py
WorkingDirectory=/home/hans/hermes-cfo
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
